# etc/cron.daily/fixtime
#!/bin/bash
# periodically update time in /var/run/fix_time
service fixtime stop

