# DEBIAN/postinst
#!/bin/bash
update-rc.d fixtime defaults
service fixtime stop
